# Library/LaunchAgents/com.vincentkoc.sync-launchd-path.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.vincentkoc.sync-launchd-path</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>launchctl setenv PATH "$PATH"; printf '%s\n' "$PATH"</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/tmp/com.vincentkoc.sync-launchd-path.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/com.vincentkoc.sync-launchd-path.log</string>
</dict>
</plist>
